# debian/adsb-history-prune.timer
[Unit]
Description=Weekly ADS-B position history pruning

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
